# orchestration/netdeneb-unit-update@.path
[Unit]
Description=Reload systemd and restart '%I' when changed

[Path]
PathChanged=/srv/netdeneb/orchestration/%I

[Install]
WantedBy=paths.target
